# crontab
* * * * * flock -xn /tmp/certs-sign.lock -c "cd /opt/rainbond-cert-controller && ./bin/cert-controller >> /proc/1/fd/1 2>&1"
